# serial-getty@.service
[Service]
ExecStart=
ExecStart=-/usr/bin/agetty --autologin root -s %I 115200,38400,9600 vt102
